# dotfiles/.config/systemd/user/dictation-shortcut.service
[Unit]
Description=Dictation shortcut on key down and key up
After=whisper-server.service
Wants=whisper-server.service

[Service]
ExecStart=%h/bin/dictation-keyboard-hook KEY_F2 /dev/input/by-id/usb-Keychron_Keychron_Q10-if02-event-kbd
Restart=on-failure
RestartSec=5
StartLimitBurst=3

[Install]
WantedBy=default.target
