# dotfiles/.config/systemd/user/dictation-shortcut.service
[Unit]
Description=Dictation shortcut on key down and key up

[Service]
ExecStart=/usr/bin/python3 %h/bin/dictation-keyboard-hook KEY_F2 /dev/input/by-id/usb-Keychron_Keychron_Q10-if02-event-kbd
Restart=on-failure
Group=input

[Install]
WantedBy=default.target
