# dotfiles/.config/systemd/user/jailmcphttp.service
[Unit]
Description=Jail MCP HTTP Server
After=network.target

[Service]
Type=simple
ExecStart=/home/vacation/bin/jailmcphttp
Restart=on-failure
RestartSec=5
Environment="JAIL_MCP_TIMEOUT=15s"
Environment="JAIL_MCP_BACKGROUND_TIMEOUT=5m"
Environment="JAIL_MCP_TRANSPORT=mcp-proxy"
Environment="PATH=/home/vacation/bin:/home/vacation/.local/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
