# dotfiles/.config/systemd/user/whisper-server-cpu.service
[Unit]
Description=Whisper.cpp inference server (CPU)
# note: ffmpeg is a dependency
Conflicts=whisper-server.service

[Service]
ExecStart=/usr/bin/whisper-server \
    -m /usr/share/whisper.cpp-model-large-v3-turbo-q5_0/ggml-large-v3-turbo-q5_0.bin \
    --host 127.0.0.1 --port 8765 \
    --language auto --suppress-nst --no-timestamps --convert \
    --no-gpu --threads 6
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target