# dotfiles/.config/systemd/user/plasma-kwin_wayland.service.d/override.conf
[Service]
Environment=KWIN_SOFTWARERENDERING=1
