# etc/systemd/system/ollama-unload.service
[Unit]
Description=Unload Ollama models from VRAM before sleep
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/home/vacation/Desktop/interface/bin/ollama-unload

[Install]
WantedBy=sleep.target